# logrotate.cron
0 0 * * * /usr/sbin/logrotate /etc/suricata/logrotate.conf
