# suricata-update.cron
0 0 * * * /etc/suricata/suricata-update.sh >> /var/log/suricata/update.log 2>&1