# suricata-update.cron
#!/bin/sh
# Cron update of suricata IDS/IPS
/usr/local/bin/suricata-update
