# suricata-update.cron
0 0 * * * /bin/sh /suricata-update-script.sh >> /var/log/suricata/update.log 2>&1