# update.cron
# Cron update of suricata IDS/IPS
@daily   root    suricata-update
